# firmware/systemd/bluesignal-wqm.service
[Unit]
Description=BlueSignal WQM-1 Main Service
After=network.target bluetooth.target gpsd.service chrony.service
Wants=bluesignal-ble.service bluesignal-led.service

[Service]
Type=simple
User=bluesignal
Group=bluesignal
WorkingDirectory=/opt/bluesignal/firmware
ExecStart=/opt/bluesignal/firmware/venv/bin/python -m src.main
Restart=always
RestartSec=10
WatchdogSec=60
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/bluesignal/firmware

# Hardware access
SupplementaryGroups=gpio i2c spi dialout

# Security hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/bluesignal /var/log/bluesignal /etc/bluesignal /opt/bluesignal/config /dev/watchdog
ReadOnlyPaths=/sys/bus/w1 /sys/class/thermal /proc/cpuinfo /sys/class/net
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
